# script/ndprd
#!/bin/ksh

daemon="/usr/local/sbin/ndprd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
