# script/ndrd
#!/bin/ksh

daemon="/usr/local/sbin/ndrd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
